# etc/cron.daily/remove-old-escenic-log-files
#! /usr/bin/env bash

## Uncomment the following line to have daily removal of old log
## files.
# /etc/init.d/ece remove-old-log-files

