# etc/cron.d/periodic-check
MAILTO=root

* * * * *     nagios /bin/sleep 20; if [ -x /usr/lib/periodic-check/main ]; then /usr/lib/periodic-check/main; fi

