# etc/cron.d/periodic-check
MAILTO=root

* * * * *     nagios if [ -x /usr/lib/periodic-check/main ]; then /usr/lib/periodic-check/main; fi

